# overlay/etc/systemd/system/rpi-data-init.service
[Unit]
# pi-appliance-kit — provision /data before it is mounted and before the read-only
# overlay takes effect. Idempotent: the script exits 0 immediately once a partition
# labeled "data" exists (blkid -L data guard), so it is a fast no-op on later boots.
# (A /data/.provisioned marker cannot be used here — this runs BEFORE data.mount, so
# a marker under /data would land on the pre-mount dir and never persist.)
Description=Provision writable /data partition (first boot)
DefaultDependencies=no
# Order after udev has processed coldplug events (udev-trigger) — NOT the
# deprecated udev-settle, which waited for ALL events and cost ~1.7s on the boot
# critical path. The root disk /dev/mmcblk0 already exists (root mounts from it),
# and the script does its own partprobe + settle for any partition it creates.
After=systemd-udev-trigger.service
Before=data.mount local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/rpi-data-init.sh

[Install]
WantedBy=local-fs.target
